# cron/imapfilter-dispatcher.cron
# IMAPFilter Dispatcher — läuft jede Minute
# Installieren: cp imapfilter-dispatcher /etc/cron.d/imapfilter-dispatcher
* * * * * www-data /usr/bin/php /var/www/imapfilter-ui/cron/dispatcher.php
